# dot_config/systemd/user/fish-extra-completions.path
[Unit]
Description=Fish extra completions

[Path]
PathChanged=/snap/bin/lxc

[Install]
WantedBy=default.target
